# scripts/com.stock-dashboard.cloudflare-tunnel.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.stock-dashboard.cloudflare-tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/yanwei/.local/bin/stock-dashboard-tunnel</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/tmp/stock-dashboard-tunnel.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/stock-dashboard-tunnel.err</string>
</dict>
</plist>
